# toolstpl/scripts/platform/macos/tool.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>__PLIST_ID__</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/__BIN_NAME__</string>
    <!-- TODO: 添加你的工具所需的启动参数 -->
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>__HOME__/.config/__TOOL_NAME__/logs/__TOOL_NAME__.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.config/__TOOL_NAME__/logs/__TOOL_NAME__.log</string>
  <key>WorkingDirectory</key>
  <string>__HOME__/.config/__TOOL_NAME__</string>
</dict>
</plist>
