# smartxray/scripts/platform/openwrt/smartxray
#!/bin/sh /etc/rc.common
# SmartXray OpenWrt init script
START=99
STOP=10
USE_PROCD=1

SMARTXRAY_USER="${SMARTXRAY_USER:-root}"
DATA_DIR="$(eval echo ~$SMARTXRAY_USER)/.config/smartxray"
CONF_FILE="$DATA_DIR/config.json"

start_service() {
    # 确保数据目录存在
    mkdir -p "$DATA_DIR/logs"

    procd_open_instance
    procd_set_param command /usr/local/bin/xray run -c "$CONF_FILE"
    procd_set_param respawn
    procd_set_param stdout 1
    procd_set_param stderr 1
    procd_set_param user "$SMARTXRAY_USER"
    procd_close_instance
}
